# rc.local
#!/bin/bash
#script to raise multiple wireguard interfaces at boot

list=$(ls -1 /etc/wireguard/*.conf)

for nic in ${list[@]}
	do
	wg-quick up $nic
	sleep 0.3
done
echo module wireguard +p > /sys/kernel/debug/dynamic_debug/control
